# scripts/deploy/server-config/systemd/clickhouse-server.service.d/override.conf
[Service]
Environment=CLICKHOUSE_WATCHDOG_ENABLE=0
