# scripts/deploy/server-config/systemd/logrotate.timer.d/override.conf
[Timer]
OnCalendar=
OnCalendar=hourly
AccuracySec=15m
Persistent=true
